# scorchd.service
[Unit]
Description=scorchd — BLE thermal printer daemon
After=bluetooth.target

[Service]
Type=simple
WorkingDirectory=%h/Repos/re_thermal_printer
ExecStart=/home/cael/.local/bin/uv run print.py --daemon 30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
