# scorchd.service
[Unit]
Description=scorchd — BLE thermal printer daemon
After=bluetooth.target

[Service]
Type=simple
ExecStart=%h/.local/bin/scorchd --daemon 30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
